# root/etc/services.d/freshclam/run
#!/usr/bin/with-contenv bash

exec freshclam -d >/dev/null