# root/etc/services.d/fail2ban/run
#!/usr/bin/with-contenv bash

if [[ ${FAIL2BAN_ENABLED} != "disabled" ]]; then
    exec fail2ban-client -x -f start
fi
